# files/master/etc/uci-defaults/90-adguardhome
uci batch << EOF

set adguardhome.config.enabled='1'
set adguardhome.config.config='/etc/adguardhome/adguardhome.yaml'

commit
EOF