# files/master/etc/uci-defaults/90-cloudflared
uci batch << EOF

set cloudflared.config.enabled='$IS_ENABLE'
set cloudflared.config.token='$CLOUDFLARED_TUNNEL_TOKEN'

commit
EOF