# files/master/etc/uci-defaults/90-cloudflared
uci batch << EOF

set cloudflared.config.enabled='1'
set cloudflared.config.token='xxoo'

commit
EOF