# files/master/etc/uci-defaults/90-gecoosac
uci batch << EOF

set gecoosac.config.enabled='$IS_ENABLE'
set gecoosac.config.log='1'

commit
EOF