# files/master/etc/uci-defaults/90-nikki-feed
echo "src/gz nikki $NIKKI_FEED_URL" >> /etc/opkg/customfeeds.conf

opkg-key add /etc/nikki/profiles/key-build.pub; sleep 1;

rm -rf /etc/nikki/profiles/key-build.pub