# files/master/etc/uci-defaults/90-nikki-feed
curl -s -L https://github.com/nikkinikki-org/OpenWrt-nikki/raw/refs/heads/main/feed.sh | ash
sleep 1;