# files/master/etc/uci-defaults/90-shairport-sync
uci batch << EOF

set shairport-sync.shairport_sync.enabled='$IS_ENABLE'
set shairport-sync.shairport_sync.name='罗家音响'

commit
EOF