# files/master/etc/uci-defaults/91-system
uci batch << EOF

set luci.main.lang='zh_cn'
set system.@system[0].hostname='LohomeWrt'
set system.@system[0].zonename='Asia/Shanghai'
set system.@system[0].timezone='CST-8'
set system.@system[0].cronloglevel='8'
set system.@system[0].conloglevel='7'

del system.ntp.server
add_list system.ntp.server='ntp.tencent.com'
add_list system.ntp.server='ntp.aliyun.com'
add_list system.ntp.server='ntp.ntsc.ac.cn'
add_list system.ntp.server='cn.ntp.org.cn'
add_list system.ntp.server='time1.apple.com'

commit
EOF