# files/master/etc/uci-defaults/99-custom
uci batch << EOF

set network.lan.ipaddr='10.0.50.2'

commit
EOF