# files/master/etc/uci-defaults/99-dhcp
uci batch << EOF

del dhcp.lan.ra_slaac
del dhcp.lan.dhcpv6
set dhcp.lan.limit='200'
set dhcp.lan.leasetime='6h'
add_list dhcp.lan.dhcp_option='3,192.168.50.1'
add_list dhcp.lan.dhcp_option='6,192.168.50.1'
del dhcp.lan.ra_flags
add_list dhcp.lan.ra_flags='none'
set dhcp.lan.dns_service='0'

commit
EOF