# files/master/etc/uci-defaults/99-network
uci batch << EOF

del network.globals.ula_prefix

set network.lan.ipaddr='$NETWORK_LAN_IP'
set network.lan.netmask='255.255.255.0'
set network.lan.ip6assign='64'
set network.lan.ip6ifaceid='eui64'

uci set network.wan.proto='pppoe'
uci set network.wan.ipv6='auto'
uci set network.wan.username='$NETWORK_WAN_USERNAME'
uci set network.wan.password='$NETWORK_WAN_PASSWORD'

set network.MODEM=interface
set network.MODEM.proto='static'
set network.MODEM.device='eth1'
set network.MODEM.ipaddr='192.168.1.111'
set network.MODEM.gateway='192.168.1.1'
set network.MODEM.netmask='255.255.255.0'
set network.MODEM.defaultroute='0'

commit
EOF